# crontab
* * * * * /usr/local/bin/php /var/www/html/artisan ScheduleGenerateVoucher:generateVoucher